# turtle-soup/soup.ini
[program:turtle-soup]
command=/opt/cedarstar/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8012
directory=/opt/cedartoy/turtle-soup/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/turtle-soup.err.log
stdout_logfile=/var/log/turtle-soup.out.log
