# deploy/systemd/cedartoy-clean-guest-saves.service
[Unit]
Description=Cedar Toy guest save cleanup

[Service]
Type=oneshot
WorkingDirectory=/opt/cedartoy
ExecStart=/usr/bin/python3 /opt/cedartoy/scripts/clean_guest_saves.py
